# cookbooks/chef-streambot-aws/templates/default/update-instance-infos.conf.erb
description "Updates AWS instance tags on boot"

start on runlevel [2345]

task
<% @scripts.keys.sort.each do |idx| %>
bash <%= @scripts[idx] %>
<% end %>
exit 0