# cookbooks/chef-streambot-aws/templates/default/update-instance-infos.conf.erb
description "Updates AWS instance tags on boot"

start on runlevel [2345]

task
<% @scripts.each do |idx, script| %>
bash <%= script %>
<% end %>
exit 0