# cookbooks/chef-streambot-aws/templates/default/update-aws-tags.conf.erb
description "Updates AWS instance tags on boot"

start on runlevel [2345]

task
exec #{@script}
