# cookbooks/chef-streambot-api/templates/default/streambot-api.conf.erb
#!upstart
# Streambot API upstart script
# This script will start/stop a Streambot API server
description "Start and stop the Streambot API server"
version "0.1"
author "Martin Biermann"

start on [2345]
stop on [016]

respawn

setuid <%= @user %>

exec <%= @binary %> -c <%= @config %> >> <%= @log %> 2>&1
